# docker/install/conf.d/avahi.conf
[program:avahi]
command = service avahi-daemon start
user = root
autostart=true
autorestart=true
