# docker/install/conf.d/dbus.conf
[program:dbus]
command = service dbus start
user = root
autostart=true
autorestart=true
