# docker/install/conf.d/prep.conf
[program:prep]
command = bash /tmp/install/prep.sh
startsecs=0
autorestart=false
startretries=1
priority=1
