# docker/install/conf.d/ssh.conf
[program:ssh]
command = service ssh start
user = root
autostart=true
autorestart=true
